# crontab
*/10 * * * * python /SeasonRecordCompareAtGame.py >> /var/log/season-record-compare-at-game.log 2>&1
