# crontab
*/10 * * * * python /SeasonRecordCompareAtGame.py