# crontab
PYBASEBALL_CACHE=/cache
CONFIG_FILE_PATH=/config/config.yaml
SECRETS_FILE_PATH=/config/secrets.yaml
*/10 * * * * /SeasonMetricsCompareAtGame.py >> /var/log/season-record-compare-at-game.log 2>&1
